# tp/3/suggestion_correction/yt-v2/yt-v2.service
[Unit]
Description=Youtube downloader daemon

[Service]
ExecStart=/srv/yt/yt-v2.sh
User=yt
WorkingDirectory=/srv/yt

[Install]
WantedBy=multi-user.target

